# results/generated_units/satpi-pass-20260501T022342Z-METEOR-M2-4.timer
[Unit]
Description=SATPI timer for METEOR-M2 4 (south_to_north) at 2026-05-01T02:23:42Z

[Timer]
OnCalendar=2026-05-01 02:23:42 UTC
Persistent=true
Unit=satpi-pass-20260501T022342Z-METEOR-M2-4.service

[Install]
WantedBy=timers.target
