# results/generated_units/satpi-pass-20260503T113433Z-METEOR-M2-3.service
[Unit]
Description=SATPI pass receiver for METEOR-M2 3 (north_to_south) at 2026-05-03T11:34:33Z
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=andreas
WorkingDirectory=/home/andreas/satpi
ExecStart=/usr/bin/python3 /home/andreas/satpi/bin/receive_pass.py --pass-file /home/andreas/satpi/results/generated_units/satpi-pass-20260503T113433Z-METEOR-M2-3.pass.json
