# results/generated_units/satpi-pass-20260506T213830Z-METEOR-M2-3.service
[Unit]
Description=SATPI pass receiver for METEOR-M2 3 (south_to_north) at 2026-05-06T21:38:30Z
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=andreas
WorkingDirectory=/home/andreas/satpi
ExecStart=/usr/bin/python3 /home/andreas/satpi/bin/receive_pass.py --pass-file /home/andreas/satpi/results/generated_units/satpi-pass-20260506T213830Z-METEOR-M2-3.pass.json
