# results/generated_units/satpi-pass-20260503T221821Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (south_to_north) at 2026-05-03T22:18:21Z

[Timer]
OnCalendar=2026-05-03 22:18:21 UTC
Persistent=true
Unit=satpi-pass-20260503T221821Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
