# results/generated_units/satpi-pass-20260503T161150Z-METEOR-M2-4.timer
[Unit]
Description=SATPI timer for METEOR-M2 4 (north_to_south) at 2026-05-03T16:11:50Z

[Timer]
OnCalendar=2026-05-03 16:11:50 UTC
Persistent=true
Unit=satpi-pass-20260503T161150Z-METEOR-M2-4.service

[Install]
WantedBy=timers.target
