# results/generated_units/satpi-pass-20260505T032315Z-METEOR-M2-4.timer
[Unit]
Description=SATPI timer for METEOR-M2 4 (south_to_north) at 2026-05-05T03:23:15Z

[Timer]
OnCalendar=2026-05-05 03:23:15 UTC
Persistent=true
Unit=satpi-pass-20260505T032315Z-METEOR-M2-4.service

[Install]
WantedBy=timers.target
