# results/generated_units/satpi-pass-20260506T033822Z-METEOR-M2-4.timer
[Unit]
Description=SATPI timer for METEOR-M2 4 (south_to_north) at 2026-05-06T03:38:22Z

[Timer]
OnCalendar=2026-05-06 03:38:22 UTC
Persistent=true
Unit=satpi-pass-20260506T033822Z-METEOR-M2-4.service

[Install]
WantedBy=timers.target
