# results/generated_units/satpi-pass-20260430T104032Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (north_to_south) at 2026-04-30T10:40:32Z

[Timer]
OnCalendar=2026-04-30 10:40:32 UTC
Persistent=true
Unit=satpi-pass-20260430T104032Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
