# results/generated_units/satpi-pass-20260502T155652Z-METEOR-M2-4.timer
[Unit]
Description=SATPI timer for METEOR-M2 4 (north_to_south) at 2026-05-02T15:56:52Z

[Timer]
OnCalendar=2026-05-02 15:56:52 UTC
Persistent=true
Unit=satpi-pass-20260502T155652Z-METEOR-M2-4.service

[Install]
WantedBy=timers.target
