# results/generated_units/satpi-pass-20260505T215143Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (south_to_north) at 2026-05-05T21:51:43Z

[Timer]
OnCalendar=2026-05-05 21:51:43 UTC
Persistent=true
Unit=satpi-pass-20260505T215143Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
