# results/generated_units/satpi-pass-20260504T220500Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (south_to_north) at 2026-05-04T22:05:00Z

[Timer]
OnCalendar=2026-05-04 22:05:00 UTC
Persistent=true
Unit=satpi-pass-20260504T220500Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
