# results/generated_units/satpi-pass-20260506T213830Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (south_to_north) at 2026-05-06T21:38:30Z

[Timer]
OnCalendar=2026-05-06 21:38:30 UTC
Persistent=true
Unit=satpi-pass-20260506T213830Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
