# results/generated_units/satpi-pass-20260503T113433Z-METEOR-M2-3.timer
[Unit]
Description=SATPI timer for METEOR-M2 3 (north_to_south) at 2026-05-03T11:34:33Z

[Timer]
OnCalendar=2026-05-03 11:34:33 UTC
Persistent=true
Unit=satpi-pass-20260503T113433Z-METEOR-M2-3.service

[Install]
WantedBy=timers.target
